# packages/qna-poc/deploy/proj-k-agent.service
# Project K QnA — systemd 서비스 파일 (FastAPI + uvicorn)
# 위치: /etc/systemd/system/proj-k-agent.service

[Unit]
Description=Project K QnA Agent (FastAPI)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/proj-k-agent/packages/qna-poc
EnvironmentFile=/home/ubuntu/proj-k-agent/packages/qna-poc/.env
ExecStart=/home/ubuntu/proj-k-agent/.venv/bin/python -m uvicorn src.api:app --host 127.0.0.1 --port 8088
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
